# daemon-cli/firmware-code/profiles/rc_car_pi_arduino/raspberry_pi/systemd/daemon-claw.service
[Unit]
Description=DAEMON RC Car Claw Node (arm)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=treehacks
WorkingDirectory=/home/treehacks/Daemon
ExecStart=/usr/bin/python3 /home/treehacks/Daemon/daemon-cli/firmware-code/profiles/rc_car_pi_arduino/raspberry_pi/claw_daemon_node.py --gpio 18 --port 8767 --node-id arm --pin-factory lgpio --detach-after-move
Restart=always
RestartSec=1
StandardOutput=append:/home/treehacks/claw_node.log
StandardError=append:/home/treehacks/claw_node.log

[Install]
WantedBy=multi-user.target
